# dot_config/systemd/user/change-wallpaper.service
[Unit]
Description=Change Hyprland Wallpaper

[Service]
Type=oneshot
ExecStart=%h/bin/change_wallpaper.sh
