# dot_config/systemd/user/change-wallpaper.timer
[Unit]
Description=Run wallpaper changer hourly

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
