# dot_config/systemd/user/update-calendars.timer
[Unit]
Description=Run Calendar updater hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
